# target/linux/ralink/base-files/etc/hotplug.d/wifi/001-wifi.sh
#!/bin/sh
logger -t wifi   "WIFIDEV=$WIFIDEV  ACTION=$ACTION  MAC=$MAC  SPEED=$SPEED SIGNAL=$SIGNAL"

